# scripts/systemd/seedbuster-fingerprint-refresh.timer
[Unit]
Description=Weekly SeedBuster visual fingerprint refresh

[Timer]
OnCalendar=Sun *-*-* 03:30:00
Persistent=true
RandomizedDelaySec=1800

[Install]
WantedBy=timers.target
